# .crontab
0,59/15 * * * * offlineimap
13 * * * * ps aux|grep offlineimap|grep python|awk '{print $2}'|xargs kill

