# .crontab
0-59/5 * * * * offlineimap
0 0 * * * ~/bin/linux/cron-code-reading
0 * * * * ~/bin/linux/backup_where

