# .crontab
0-59/15 *  * * * . ~/.bash_path; offlineimap
0       0  * * * . ~/.bash_path; ~/bin/linux/cron-code-reading
0       *  * * * . ~/.bash_path;  ~/bin/linux/backup_where
0       *  * * * . ~/.bash_path;  ~/bin/linux/download-google-sdk.sh
0       *  * * * . ~/.bash_path; ~/bin/linux/jira-getissues
30      22 * * * . ~/.bash_path; ~/bin/linux/update-repo-bhj
0-59/15 *  * * * . ~/.bash_path; ~/bin/git-mac-branch ~/windows-config/
0-59/15 *  * * * . ~/.bash_path; ~/bin/notify-lost-offlineimap
0       20 * * * . ~/.bash_path; ~/bin/git-to-gmail ~/notes
1       1  * * * . ~/.bash_path; ~/bin/linux/download-today
0       0  * * * . ~/.bash_path; ~/bin/linux/update-repo-projects; ~/bin/linux/cron-code-reading
0       20 * * * . ~/.bash_path; back-to-gmail ~/.ecompleterc ~/.bash_history ~/.bash_path
11      20 * * * . ~/.bash_path; cd ~/cygwin && mirrorcyg.pl | mf-download.pl -n 3

