# .crontab
0,59/15 * * * * offlineimap
13 * * * * ps aux|grep offlineimap|grep python|awk '{print $2}'|xargs kill
0 23 * * * cd ~/windows-config && git pull && git commit -a -m 'crontab commit' && git push
