# .crontab
0-59/15 *  * * * . ~/.bash_path; ~/bin/bhj-cron-runner 15m 15m
0       22 * * * . ~/.bash_path; ~/bin/bhj-cron-runner daily 10h
0       *  * * * . ~/.bash_path; ~/bin/bhj-cron-runner hourly 1h
0       23 * * 4 . ~/.bash_path; ~/bin/bhj-cron-runner weekly 10h

