# .crontab
0-59/15 * * * * offlineimap
0 0 * * * ~/bin/linux/cron-code-reading
0 * * * * ~/bin/linux/backup_where
0 * * * * ~/bin/linux/jira-getissues
30 22 * * * bash -c ". ~/.bashrc; ~/bin/linux/update-repo-bhj"
0 * * * * ~/bin/config-save-to-gmail
