# work_bin/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
#/usr/src/redhat/SOURCES/ypbind.init start
/etc/init.d/ypbind start
(/etc/init.d/netfs restart;/etc/init.d/clearcase start )

date >/home/a22242/date.txt
