# etc/subdir-symlinks/.config/autostart/xbindkeys.desktop
[Desktop Entry]
Exec=env start-xbindkeys
Name=Start Xbindkeys
Terminal=false
Type=Application
