# etc/udev/rules.d/90-bhj-personal.rules
SUBSYSTEM=="usb", ATTR{idVendor}=="096e", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="04e8", MODE="0600", OWNER="bhj", RUN+="/home/bhj/bin/ask-to-sync-org"
SUBSYSTEM=="usb", ATTR{idVendor}=="0955", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="0bb4", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="1949", MODE="0600", OWNER="bhj", RUN+="/home/bhj/bin/ask-to-sync-org"
SUBSYSTEM=="usb", ATTR{idVendor}=="18d1", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="05c6", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="b155", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="045e", RUN+="/home/bhj/bin/re-xmodmap"
SUBSYSTEM=="usb", ATTR{idVendor}=="29a9", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="1004", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="2717", MODE="0600", OWNER="bhj"
SUBSYSTEM=="usb", ATTR{idVendor}=="05ac", MODE="0600", OWNER="bhj"
